# bundled/ponytail/src/gnome-ponytail-daemon.service.in
[Unit]
Description=GNOME Ponytail daemon

[Service]
Type=dbus
BusName=org.gnome.Ponytail
ExecStart=@libexecdir@/gnome-ponytail-daemon
Restart=on-failure
