# bundled-bin/ppc64le/usr/lib/systemd/user/gnome-ponytail-daemon.service
[Unit]
Description=GNOME Ponytail daemon

[Service]
Type=dbus
BusName=org.gnome.Ponytail
ExecStart=/usr/libexec/gnome-ponytail-daemon
Restart=on-failure
