# etc/systemd/system/anabot.service
[Unit]
Description=Anaconda 
After=anaconda.service
Wants=anabot-prepare.service
After=anabot-prepare.service

[Service]
ExecStart=/opt/start-anabot.py
Environment="PYTHONPATH=/opt/lib/python2.7/site-packages"

[Install]
WantedBy=anaconda.service
